# iso-profiles/garuda/gnome/desktop-overlay/etc/systemd/system/scx_lavd.service
[Unit]
# sched_ext LAVD scheduler (the Bazzite/PikaOS/CachyOS-Handheld scheduler win).
# LAVD = Latency-criticality Aware Virtual Deadline; it is dual-CCD / 3D-V-Cache aware,
# so it places latency-critical game threads on the right CCD of the 9950X3D.
Description=SCX LAVD scheduler (low-latency, X3D/CCD-aware)
Documentation=https://github.com/sched-ext/scx
After=multi-user.target

[Service]
ExecStart=/usr/bin/scx_lavd --performance
Restart=always
RestartSec=2
# If sched_ext is unavailable on a fallback kernel, don't wedge boot.
SuccessExitStatus=0 143

[Install]
WantedBy=multi-user.target
