# iso-profiles/garuda/gnome/desktop-overlay/etc/skel/.config/autostart/claw-code-setup.desktop
[Desktop Entry]
Type=Application
Name=Claw Code Setup
Comment=Provision gajae-code (Red-claw) on first login
Exec=bash -lc "$HOME/.local/bin/claw-code-setup.sh"
X-GNOME-Autostart-enabled=true
NoDisplay=true
Terminal=false
